# deploy/systemd/quasar-claw.service
[Unit]
Description=Quasar Claw Lab secure runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=quasar
Group=quasar
WorkingDirectory=/opt/quasar-claw
ExecStart=/opt/quasar-claw/bin/quasar-claw
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/tmp
RuntimeDirectory=quasar-claw
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
