# deploy/systemd/vericlaw.service
[Unit]
Description=VeriClaw secure runtime
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vericlaw
Group=vericlaw
WorkingDirectory=/opt/vericlaw
ExecStart=/opt/vericlaw/bin/vericlaw
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/tmp
RuntimeDirectory=vericlaw
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
